# crontab
* * * * * root python3 /app/app.py >> /app/crypto_data.txt
